# templates/abbr.service.j2
[Unit]
Description=abbr server
After=network.target

[Service]
Type=simple
User=abbr
ExecStart={{ abbr_home }}/bin/gunicorn_abbr

[Install]
WantedBy=multi-user.target
